# files/openhab.service
[Unit]
Description=openHAB 2
Requires=docker.service
After=docker.service

[Service]
Restart=always

WorkingDirectory=/opt/openhab/docker

# Compose up
ExecStart=/usr/bin/docker-compose up

# Compose down, remove containers and volumes
ExecStop=/usr/bin/docker-compose down

[Install]
WantedBy=multi-user.target
